# ema.timer
[Unit]
Description=Runs HelpDesk

[Timer]
# Run after booting one minute
OnBootSec=2min
Unit=ema.service

[Install]
WantedBy=multi-user.target