# ema.service
[Unit]
Description=helpdesk

[Service]
Type=simple
ExecStart=/root/emahlpdsk.sh